# build/conf/crontab
####################################
# CRONTAB
#
# [utilisateur] : utilisateur avec lequel on exécute la commande dans le conteneur (normalement www-data)
# [commande] : commande à exécuter
#
# Après modification ne pas oublier de lancer "make rebuild"
#
#####################################

#* * * * * [utilisateur] [commande] >> /dev/null 2>&1

# Laravel Schedule
* * * * * www-data cd /var/www/app/ && php artisan schedule:run >> /dev/null 2>&1
