# build/conf/supervisord.conf
[supervisord]
nodaemon=true
pidfile = /root/supervisord.pid
logfile = /var/log/supervisord.log
logfile_maxbytes = 50MB
logfile_backups=10

[program:cron]
command=cron -f
autostart=true
autorestart=true

[program:apache2]
command=/usr/bin/init
autostart=true
autorestart=true
